# tools/actions-runner/evidence-producer.service
# SPDX-FileCopyrightText: Copyright (c) 2026 contributors to SOL ExecBench ROCm Port
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=SOL ExecBench evidence producer GitHub Actions runner
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment=EVIDENCE_SOURCE_ROOT=%h/.local/share/sol-execbench/evidence-input
WorkingDirectory=%h/.local/share/actions-runner/SOL-ExecBench-ROCm
ExecStart=%h/.local/share/actions-runner/SOL-ExecBench-ROCm/run.sh
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=5min

[Install]
WantedBy=default.target
